# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:web_healthcheck]
command=python -c "import http.server, socketserver, os; port=int(os.environ.get('PORT', 10000)); socketserver.TCPServer.allow_reuse_address=True; httpd=socketserver.TCPServer(('0.0.0.0', port), http.server.SimpleHTTPRequestHandler); print(f'Healthcheck HTTP server running on port {port}'); httpd.serve_forever()"
autorestart=true
priority=1

[program:mosquitto]
command=mosquitto -c /etc/mosquitto/mosquitto.conf
autorestart=true
priority=2

[program:boiler_engine]
command=python engine/boiler_engine.py
autorestart=true
priority=3

[program:anomaly_detector]
command=python engine/anomaly_detector.py
autorestart=true
priority=4

[program:historian_service]
command=python engine/historian_service.py
autorestart=true
priority=4
